# bt_bind1@.service
[Unit]
Desription=RFCOMM aut-bind 1
After=bluetooth.target

[Service]
ExecStart=/bin/bash /home/pi/rfcomm_autobind/bt_bind.sh %i 1
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target

